# features/chost/file.include/etc/profile.d/65-chost.sh
#!/usr/bin/env bash

# Source bash completion the k8s-tools
source <(kubeadm completion bash)
source <(kubectl completion bash)
source <(helm completion bash)